# .config/systemd/user/check-battery.service
[Unit]
Description=Check Battery Level

[Install]
WantedBy=default.target

[Service]
ExecStart=/home/nico/.scripts/xctl check-battery 20

