# .scripts/feeds/jobs.crontab
*/5 * * * * make -f /home/nico/.scripts/feeds/makefile > ~/.local/share/feeds/cron.log 2>&1
